# server/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:tibia]
startsecs=0
autorestart=false
command=/opt/startServer.sh

